# systemd/ai-chat-server.service
[Unit]
Description=AI MCP Chat Server for operation-dbus
Documentation=https://github.com/operation-dbus
After=network.target

[Service]
Type=simple
User=jeremy
WorkingDirectory=/home/jeremy/op-dbus-staging
Environment="NODE_ENV=production"
Environment="PATH=/home/jeremy/.nvm/versions/node/v22.19.0/bin:/usr/local/bin:/usr/bin:/bin"
Environment="AI_PROVIDER=huggingface"
Environment="HF_EMAIL=jeremy@ghostbridge.tech"
Environment="HF_TOKEN=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"
Environment="HF_MODEL=Qwen/Qwen2.5-VL-7B-Instruct"
ExecStart=/home/jeremy/.nvm/versions/node/v22.19.0/bin/node /home/jeremy/op-dbus-staging/chat-server.js
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ai-chat

[Install]
WantedBy=multi-user.target
