# systemd/mcp-chat-server.service
[Unit]
Description=MCP Chat Server with Unified Introspection
Documentation=https://github.com/operation-dbus
After=network.target dbus.service
Wants=dbus.service

[Service]
Type=simple
User=jeremy
WorkingDirectory=/home/jeremy/op-dbus-staging
ExecStart=/home/jeremy/.cargo/bin/mcp-chat
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=mcp-chat

# Environment variables for AI integration
Environment="RUST_LOG=info"
Environment="RUST_BACKTRACE=1"

# AI Configuration - REQUIRED (AI is the brain of the system)
# Get your API key from: https://ollama.com
Environment="OLLAMA_API_KEY=your-api-key-here"
Environment="OLLAMA_DEFAULT_MODEL=qwen"

# Resource limits
MemoryLimit=2G
CPUQuota=80%

[Install]
WantedBy=multi-user.target
